# com.gaganpreet.citizenshiptracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gaganpreet.citizenshiptracker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/gaganpreet/gemini/personal_projects/citizenship-tracker/.venv/bin/python3</string>
        <string>/Users/gaganpreet/gemini/personal_projects/citizenship-tracker/scraper.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/gaganpreet/gemini/personal_projects/citizenship-tracker</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/gaganpreet/gemini/personal_projects/citizenship-tracker/tracker_local.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gaganpreet/gemini/personal_projects/citizenship-tracker/tracker_error.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
